# cron/icvpn-meta-zone.cron
# Rebuild icvpn zones once per hour
0 * * * * $HOME/icvpn-aux/scripts/update-icvpn-zone.sh
